# debian/tibber/DEBIAN/postinst
#!/bin/bash
systemctl enable tibber.service
systemctl start tibber.service
