# debian/tibber/DEBIAN/postinst
#!/bin/bash
cd /opt/tibber
if [ -f "appsettings.user" ]; then
  mv appsettings.user appsettings.json
fi

systemctl enable tibber.service
systemctl start tibber.service
