# debian/tibber/DEBIAN/preinst
#!/bin/bash
systemctl stop tibber.service || true
systemctl disable tibber.service || true
